# root/etc/s6-overlay/s6-rc.d/adguardhome/run
#!/usr/bin/with-contenv bash

s6-svwait -U /var/run/service/unbound

cmd="AdGuardHome --work-dir /config/adguardhome --host 0.0.0.0"
echo "*** executing => $cmd"
exec $cmd
